# deployment/systemd/prosody-cert-monitor.timer
[Unit]
Description=Run Prosody Certificate Monitor
Documentation=https://github.com/your-org/prosody-xmpp-server
Requires=prosody-cert-monitor.service

[Timer]
# Run every hour
OnCalendar=hourly
# Run 5 minutes after boot
OnBootSec=5min
# Randomize start time by up to 10 minutes to avoid load spikes
RandomizedDelaySec=600
# Ensure timer is persistent across reboots
Persistent=true

[Install]
WantedBy=timers.target